# data/com.kaira.cron.test.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.EN">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kaira.cron.test</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/derin/Desktop/ANTIGRAVITY-AGENT/anti-claw/venv/bin/python</string>
        <string>/Users/derin/Desktop/ANTIGRAVITY-AGENT/anti-claw/data/send_cron_test.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>10</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/derin/Desktop/ANTIGRAVITY-AGENT/anti-claw/data/logs/cron_test.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/derin/Desktop/ANTIGRAVITY-AGENT/anti-claw/data/logs/cron_test.err</string>
</dict>
</plist>
